# vagrant_scripts/ipfsd.service
[Unit]
Description=IPFS daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/ipfs daemon

[Install]
WantedBy=multiuser.target